# docker/selfhost/crontab
0 6 * * * /app/docker/selfhost/run-sync.sh
